# scripts/sns-poster.service
[Unit]
Description=SNS Poster API Service - Multi-platform Social Media Publishing
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
# Working directory
WorkingDirectory=/opt/sns-poster
# Binary path
ExecStart=/bin/bash -c "/opt/sns-poster/sns-poster >> /var/logs/sns-poster/sns-poster.log  2>&1"
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
Alias=sns-poster.service
